# luci-app-kubo/root/etc/uci-defaults/luci-app-kubo
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@kubo[-1]
	add ucitrack kubo
	set ucitrack.@kubo[-1].init=kubo
	commit ucitrack
EOF

rm -rf /tmp/luci-indexcache*
exit 0
